# ROS2_Pi/systemd_services/xbox360_pad_bridge.service
[Unit]
Description=Xbox Controller to UART Bridge
# No longer requires xboxdrv to be running at start - will wait for controller
After=network.target

[Service]
# Run as your user, as the script is located in your home directory
User=blankmcu
Environment="HOME=/home/blankmcu"
Type=simple
# Execute the Python script
ExecStart=/home/blankmcu/ros2_ws/controller_uart_bridge/controller_uart_bridge.py
Restart=always
RestartSec=5
# Remove restart limits to allow unlimited retries when controller connects/disconnects
StartLimitIntervalSec=0
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target